# capi.service
[Unit]
Description=CEC HTTP Bridge Service
Documentation=https://github.com/LukasParke/capi
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
User=capi
Group=capi
SupplementaryGroups=dialout
WorkingDirectory=/opt/capi

# Main service command
ExecStart=/opt/capi/capi -bind :8080 -name "CEC HTTP Bridge"

# Restart policy
Restart=on-failure
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/dev /opt/capi

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=capi

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
